# autostart/auto-theme-switch.desktop
[Desktop Entry]
Name=Auto theme switch
Comment=Desktop auto light and dark theme switch
Exec=auto-theme-switch
Type=Application
